# gpu-numa-optimizer.service
[Unit]
Description=GPU NUMA Node Affinity Optimizer
After=network.target nvidia-persistenced.service graphical.target
Wants=nvidia-persistenced.service

[Service]
Type=simple
# Standard options:
# -d: daemon mode (required for service)
# -p: physical cores only
# -s: strict memory policy
# -a: optimize all GPU processes (not just games)
ExecStart=/usr/local/bin/gpu_optimize.sh -d
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
