# gpu-numa-optimizer.service
[Unit]
Description=GPU NUMA Node Affinity Optimizer
After=network.target graphical.target
Documentation=https://github.com/mattkenn4545/gpu-numa-tune

[Service]
Type=simple
# Standard options:
# -d: daemon mode (required for service)
# -p: physical cores only
# -s: strict memory policy
# -a: optimize all GPU processes (not just games)
ExecStart=/usr/local/bin/gpu_optimize.sh -d
Restart=always
RestartSec=10
User=root

# Ensure the root management process has time to restore system settings on stop.
# SIGTERM is sent to all processes in the control group.
KillMode=control-group
TimeoutStopSec=30

# Logging to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
